# deploy/agent/probectl-ebpf-agent.service
# probectl-ebpf-agent — hardened systemd unit (U-052).
#
# Capability model (see deploy/agent/README.md for the kernel matrix):
#   kernel >= 5.8: CAP_BPF + CAP_PERFMON (the minimal pair: program/map
#                  loading + uprobe/tracepoint attachment) — preferred.
#   kernel <  5.8: replace both with CAP_SYS_ADMIN (the legacy catch-all).
# The agent is observe-only either way (the CI gate forbids enforcing
# program types); it never needs CAP_NET_ADMIN, CAP_SYS_PTRACE or root.
[Unit]
Description=probectl eBPF host agent (observe-only)
Documentation=https://github.com/imfeelingtheagi/probectl/blob/main/docs/ebpf-agent.md
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/probectl-ebpf-agent --config /etc/probectl/ebpf-agent.yaml
Restart=on-failure
RestartSec=5s

# --- privilege: minimal capabilities, no root ---
User=probectl-agent
Group=probectl-agent
AmbientCapabilities=CAP_BPF CAP_PERFMON
CapabilityBoundingSet=CAP_BPF CAP_PERFMON
NoNewPrivileges=yes
# BPF maps/ringbuf need locked memory.
LimitMEMLOCK=infinity

# --- resource caps (EBPF-007): mirror the Helm chart's limits ---
# Helm: requests 100m/128Mi, limits 1 cpu / 512Mi. A host agent must never
# starve the workloads it observes.
CPUQuota=100%
CPUWeight=100
MemoryHigh=384M
MemoryMax=512M
TasksMax=128
IOWeight=100

# --- syscall filter (mirror of deploy/agent/seccomp.json) ---
SystemCallArchitectures=native
SystemCallFilter=@system-service bpf perf_event_open
SystemCallFilter=~@mount @module @reboot @swap @obsolete @cpu-emulation @keyring ptrace
SystemCallErrorNumber=EPERM

# --- filesystem / namespace hardening ---
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/sys/kernel/btf
ReadWritePaths=/var/lib/probectl
PrivateTmp=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
